# ld.so.preload
/opt/ngx-lua-request-time/ngx_lua_request_time.so
